# package/ascend/net/gaethmon/files/gaethmon.init
#!/bin/sh /etc/rc.common

# Copyright (C) wynter

START=99
STOP=85

start() {
	echo "Start gaethmon"
	/bin/gaethmon > /dev/null 2>&1 &
}

stop() {
	echo "Stop gaethmon"
	/usr/bin/killall gaethmon
}
